# .systemd/prod/aivc-workers-study.service
[Unit]
Description=AI Church — Celery study worker (AI Bible Study multi-agent discussions)
After=redis-server.service aivc-workers.service

[Service]
User=simon
Group=www-data
WorkingDirectory=/opt/ai-church/workers
# AI Bible Study rounds are I/O-bound (OpenRouter chat calls) with no local model, so
# they get their own queue/process and never block the worship pipeline. The bridge
# (aivc-bridge.service) BLPOPs Laravel's ai:study list and dispatches tasks.study_discuss
# here. Concurrency 2 lets a couple of discussions run in parallel without overload.
ExecStart=/usr/bin/env bash -lc 'set -a; . ./.env; set +a; exec .venv/bin/celery -A tasks.celery_app worker -Q ai:study -c 2 -n study@%%h'
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
