# .systemd/prod/aivc-queue.service
[Unit]
Description=AI Church — Laravel queue worker
After=redis-server.service mysql.service

[Service]
User=deploy
Group=www-data
WorkingDirectory=/opt/ai-church/backend
ExecStart=/usr/bin/php artisan queue:work --tries=3 --sleep=1
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
